# portainer/systemd/docker-portainer-reload.service
[Unit]
Description=Refresh Portainer image and update container

[Service]
Type=oneshot

ExecStart=/bin/systemctl reload docker-portainer.service
