# nextcloud/systemd/docker-nextcloud-reload.service
[Unit]
Description=Refresh Nextcloud image and update container

[Service]
Type=oneshot

ExecStart=/bin/systemctl reload docker-nextcloud.service
