# guacamole/systemd/docker-guacamole-reload.service
[Unit]
Description=Refresh Guacamole image and update container

[Service]
Type=oneshot

ExecStart=/bin/systemctl reload docker-guacamole.service
