# grafana/systemd/docker-grafana-reload.service
[Unit]
Description=Refresh Grafana image and update container

[Service]
Type=oneshot

ExecStart=/bin/systemctl reload docker-grafana.service
