# uptime/systemd/docker-uptime-reload.service
[Unit]
Description=Refresh Uptime Kuma image and update container

[Service]
Type=oneshot

ExecStart=/bin/systemctl reload docker-uptime.service
