# traefik/systemd/docker-traefik-reload.service
[Unit]
Description=Refresh Traefik image and update container

[Service]
Type=oneshot

ExecStart=/bin/systemctl reload docker-traefik.service
